# root/etc/services.d/syncserver/run
#!/usr/bin/with-contenv bash

cd /opt/syncserver || exit

chown -R abc:abc /opt/syncserver

exec \
	s6-setuidgid abc /usr/local/bin/rails s
